# .devcontainer/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info

[program:exchange]
command=/bin/bash -lc "export NIX_CONFIG='experimental-features = nix-command flakes'; cd /workspace && nix develop --command ./gradlew bootRun --console=plain"
directory=/workspace
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/tmp/exchange.log
stderr_logfile=/tmp/exchange.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=0
